# monitor/k8s-monitor.cron
*/5 * * * * python /home/k8s-master/monitor/monitor_cpu.py
*/3 * * * * python /home/k8s-master/monitor/monitor_memory.py
0 */1 * * * python /home/k8s-master/monitor/monitor_disk.py
